# encoder/ompx-web-ui.service

# oMPX Web UI systemd Service Unit
# --------------------------------
# Runs the oMPX Web UI backend as a systemd service.
#
# Features:
#   - Starts after network is up
#   - Runs as 'ompx' user
#   - Sets working directory and environment
#   - Restarts on failure
#   - Installs for multi-user boot
#
# For more info, see: https://github.com/RobertC85/oMPX

[Unit]
Description=oMPX Web UI Backend
After=network.target

[Service]
Type=simple
User=ompx
WorkingDirectory=/workspaces/oMPX/encoder
Environment=OMPX_WEB_PORT=5000
ExecStart=/usr/bin/python3 /workspaces/oMPX/encoder/ompx-web-ui.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
